# Проект: ETL-процесс хранилища банковских транзакций/main.cron
0 5 * * * python3 /home/de3at/kela/py_scripts/etl.py